# rc.local
#!/bin/sh -e
#
#rc.local#
# Este script é executado no final de cada nível de execução multiusuário.
# Certifique-se de que o script "sairá 0" em caso de sucesso ou qualquer outro
# valor no erro.
#
# Para habilitar ou desabilitar este script basta alterar a execução
# bits.
#
# Por padrão este script não faz nada.
/etc/restartiptables.sh
saída 0
